# helix-init.sh
#!/bin/bash
### BEGIN INIT INFO
# Provides:          start.sh
# Required-Start:
# Required-Stop:
# X-Start-Before:
# Default-Start:     S 2 3 4 5
# Default-Stop:	     0 1 6
# X-Interactive:     true
# Short-Description: Start Helix modules
# Description:
### END INIT INFO

